# etc/init.d/one-context
#!/sbin/openrc-run

depend() {
	need localmount
	before networking iptables ip6tables
}

start() {
	local ret=1
	ebegin "Starting OpenNebula contextualization"
	eindent
	/usr/sbin/one-contextd
	ret=$?
	eoutdent
	return $ret
}

stop() {
	return 0
}

